# acspanel-client-ddns/root/etc/cron.d/acspanel-client-ddns-update
*/5 *	* * *	root /usr/local/acs-bin/ddns_update
